# roles/moltbot/templates/moltbot.service.j2
[Unit]
Description=MoltBot Gateway (self-hosted AI agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ moltbot_user }}
Group={{ moltbot_user }}
EnvironmentFile={{ moltbot_env_path }}
WorkingDirectory={{ moltbot_home }}
ExecStart={{ moltbot_cli_path }} gateway
Restart=on-failure
RestartSec=10
TimeoutStartSec=60

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
LockPersonality=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native

# Allow writes only where we expect state
ReadWritePaths={{ moltbot_state_dir }}

[Install]
WantedBy=multi-user.target

